# deploy/celery-beat.service
[Unit]
Description=Eandewigs Celery beat scheduler (eandewigs-api)
After=network.target redis-server.service postgresql.service
Wants=redis-server.service postgresql.service

[Service]
Type=simple
User=mmeso
Group=www-data
WorkingDirectory=/home/mmeso/backend
Environment=DJANGO_SETTINGS_MODULE=config.settings.prod
# DatabaseScheduler (django_celery_beat) keeps the schedule in Postgres, so no
# local schedule file is needed. --pidfile in /run avoids stale-pid issues.
RuntimeDirectory=eandewigs-beat
ExecStart=/home/mmeso/backend/.venv/bin/celery -A config beat \
    --loglevel=info \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler \
    --pidfile=/run/eandewigs-beat/beat.pid
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
