# deploy/celery-worker.service
[Unit]
Description=Eandewigs Celery worker (eandewigs-api)
After=network.target redis-server.service postgresql.service
Wants=redis-server.service postgresql.service

[Service]
Type=simple
User=mmeso
Group=www-data
WorkingDirectory=/home/mmeso/backend
Environment=DJANGO_SETTINGS_MODULE=config.settings.prod
# --concurrency=1 keeps memory tiny on a 512MB host; raise on a bigger box.
ExecStart=/home/mmeso/backend/.venv/bin/celery -A config worker \
    --loglevel=info \
    --concurrency=1 \
    --max-tasks-per-child=100
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
